# dist/com.user.gitchecker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Template LaunchAgent for the gitchecker server.
  Do not load this file directly — run dist/install.sh, which substitutes the
  __BIN__ and __LOG__ placeholders with absolute paths and loads it via launchctl.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.gitchecker</string>

    <key>ProgramArguments</key>
    <array>
        <string>__BIN__</string>
    </array>

    <!-- Start at login and restart on crash. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>EnvironmentVariables</key>
    <dict>
        <!-- launchd gives a minimal PATH; add Homebrew + common dirs so the
             server can find `git` (and `ssh`) when it shells out to fetch. -->
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG__</string>
    <key>StandardErrorPath</key>
    <string>__LOG__</string>
</dict>
</plist>
